# roles/duply/files/duply.cron
#!/bin/bash
duply main pre+bkp+post >>/tmp/duply.out
