# roles/duply/files/duply.cron
#!/bin/bash
duply main backup >>/tmp/duply.out
